# packages/singular-server-management/singular_server_management-0.0.11.tar.gz/singular_server_management-0.0.11/singular_server_management/file_templates/gunicorn_default.socket
[Unit]
Description=gunicorn_default socket

[Socket]
ListenStream=/run/gunicorn_default.sock

[Install]
WantedBy=sockets.target
